# mouseless/mouseless.service
[Unit]
Description=mouseless

[Service]
ExecStart=/usr/local/bin/mouseless

[Install]
WantedBy=multi-user.target
